# ansible/roles/lidi_receive/templates/diode-receive.service.j2
{{ ansible_managed | comment }}
[Unit]
Description=lidi data diode (receiver)
Documentation=https://anssi-fr.github.io/lidi/
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ lidi_user }}
Group={{ lidi_group }}
ExecStart={{ lidi_install_prefix }}/diode-receive \
    --from {{ lidi_receive_from }} \
{% if lidi_receive_to_tcp | length > 0 %}
    --to-tcp {{ lidi_receive_to_tcp }} \
{% endif %}
{% if lidi_receive_to_unix | length > 0 %}
    --to-unix {{ lidi_receive_to_unix }} \
{% endif %}
{% if lidi_mtu is not none %}
    --from-mtu {{ lidi_mtu }} \
{% endif %}
{% if lidi_block_bytes is not none %}
    --block {{ lidi_block_bytes }} \
{% endif %}
{% if lidi_repair_pct is not none %}
    --repair {{ lidi_repair_pct }} \
{% endif %}
{% if lidi_batch is not none %}
    --batch {{ lidi_batch }} \
{% endif %}
{% if lidi_receive_decode_threads | length > 0 %}
    --decode-threads {{ lidi_receive_decode_threads }} \
{% endif %}
{% if lidi_receive_cpu_affinity %}
    --cpu-affinity \
{% endif %}
{% if lidi_receive_heartbeat | length > 0 %}
    --heartbeat {{ lidi_receive_heartbeat }} \
{% endif %}
{% if lidi_receive_min_repair | length > 0 %}
    --min-repair {{ lidi_receive_min_repair }} \
{% endif %}
{% if lidi_receive_reset_timeout | length > 0 %}
    --reset-timeout {{ lidi_receive_reset_timeout }} \
{% endif %}
{% if lidi_receive_abort_timeout | length > 0 %}
    --abort-timeout {{ lidi_receive_abort_timeout }} \
{% endif %}
{% if lidi_receive_hash %}
    --hash \
{% endif %}
{% if lidi_receive_flush %}
    --flush \
{% endif %}
{% if lidi_log_file | length > 0 %}
    --log-file {{ lidi_log_file }} \
{% endif %}
{% if lidi_http_addr | length > 0 %}
    --http-addr {{ lidi_http_addr }} \
{% endif %}
    --log-level {{ lidi_log_level }}

Restart=on-failure
RestartSec=5s

# --- Hardening -------------------------------------------------------------
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
{% if lidi_receive_to_unix | length > 0 %}
# Unix socket lives under {{ lidi_receive_to_unix | dirname }} — must be writable.
ReadWritePaths={{ lidi_receive_to_unix | dirname }}
{% endif %}

[Install]
WantedBy=multi-user.target
